# HW6/custom_daemon.service
[Unit]
Description=custom_daemon

[Service]
Type=oneshot
ExecStart=/bin/bash /home/custom_daemon.sh
