# HW6/custom_daemon.timer
[Unit]
Description=custom_daemon

[Timer]
OnUnitActiveSec=30s
OnBootSec=30s

[Install]
WantedBy=timers.target
