# HW20/roles/hw202/templates/hw21.timer.j2
[Unit]
Description=hw21

[Timer]
OnUnitActiveSec= {{ deamon_var.OnUnitActiveSec }}
OnBootSec= {{ deamon_var.OnBootSec }}

[Install]
WantedBy=timers.target
