# HW20/roles/hw202/templates/hw21.service.j2
[Unit]
Description=hw21

[Service]
Type=oneshot
ExecStart=/bin/bash {{ deamon_var.ExecStart }}
